# systemd/ralph@.service
[Unit]
Description=Ralph — Autonomous coding agent (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ralph
Group=ralph

# Working directory is the target project repo
WorkingDirectory=/home/ralph/projects/%i
ExecStart=/usr/local/bin/bun run /opt/ralph/src/cli.ts start

# Per-instance environment
EnvironmentFile=/etc/ralph/%i.env

# Restart on failure with 10s delay
Restart=on-failure
RestartSec=10

# Graceful shutdown — SIGTERM triggers the shutdown handler
KillSignal=SIGTERM
TimeoutStopSec=120

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ralph-%i

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ralph/projects/%i /home/ralph/.ssh /home/ralph/.config

[Install]
WantedBy=multi-user.target
